# palace-daemon.service
[Unit]
Description=palace-daemon — mempalace HTTP/MCP gateway (system unit)
Documentation=https://github.com/techempower-org/palace-daemon
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=jp
Group=jp
WorkingDirectory=/home/jp/Projects/palace-daemon
ExecStartPre=-/usr/bin/fuser -k 8085/tcp
ExecStartPre=-/bin/rm -f /home/jp/.cache/palace-daemon/daemon-8085.lock
# Wait up to 60s for postgres in Docker container
ExecStartPre=/bin/bash -c 'until docker exec mempalace-db pg_isready -U palace -q 2>/dev/null; do sleep 2; done'
ExecStart=/home/jp/.local/share/palace-daemon/venv/bin/python main.py --palace /srv/mempalace-data/palace
Environment=PALACE_PORT=8085
Environment=PALACE_MAX_CONCURRENCY=4
# CRITICAL: forces in-process path, prevents recursive self-call deadlock.
# See palace-daemon#11.
Environment=PALACE_DAEMON_STRICT=0
# Pin the /search reranker model explicitly so deployed ranking is
# byte-stable across deploys (the code default is this same value, but
# pinning it here + the exact flashrank== in requirements.txt removes the
# "fresh deploy pulls a different model and silently reorders top-5" risk).
# See docs/evals/2026-05-30-retrieval-determinism.md.
Environment=PALACE_RERANK_MODEL=ms-marco-TinyBERT-L-2-v2
EnvironmentFile=/home/jp/.config/palace-daemon/env
Restart=always
RestartSec=5
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal
KillMode=mixed

# Lower oom_score so userland earlyoom (running on familiar at default
# 10% mem / 10% swap thresholds) doesn't pick the daemon as its
# SIGTERM victim during memory-heavy /graph or /backfill-age calls.
# The daemon's MemoryMax=2G drop-in caps its cgroup; this just keeps
# system-wide pressure from preempting it. -500 displaces the daemon
# below routine python processes without making it completely
# invulnerable. See palace-daemon#154.
OOMScoreAdjust=-500

[Install]
WantedBy=multi-user.target
