# palace-daemon.service
[Unit]
Description=palace-daemon — mempalace HTTP/MCP gateway (system unit)
Documentation=https://github.com/techempower-org/palace-daemon
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=jp
Group=jp
WorkingDirectory=/home/jp/Projects/palace-daemon
ExecStartPre=-/usr/bin/fuser -k 8085/tcp
ExecStartPre=-/bin/rm -f /home/jp/.cache/palace-daemon/daemon-8085.lock
# Wait up to 60s for postgres in Docker container
ExecStartPre=/bin/bash -c 'until docker exec mempalace-db pg_isready -U palace -q 2>/dev/null; do sleep 2; done'
ExecStart=/home/jp/.local/share/palace-daemon/venv/bin/python main.py --palace /srv/mempalace-data/palace
Environment=PALACE_PORT=8085
Environment=PALACE_MAX_CONCURRENCY=4
# CRITICAL: forces in-process path, prevents recursive self-call deadlock.
# See palace-daemon#11.
Environment=PALACE_DAEMON_STRICT=0
EnvironmentFile=/home/jp/.config/palace-daemon/env
Restart=always
RestartSec=5
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal
KillMode=mixed

[Install]
WantedBy=multi-user.target
