# scripts/syncthing-keepalive/syncthing-keepalive.service
[Unit]
Description=Restart syncthing@%i if it has clean-exited
Documentation=https://github.com/techempower-org/palace-daemon/issues/92
After=network-online.target

[Service]
Type=oneshot
# `is-active` returns 0 when the unit is running and non-zero otherwise; we
# only fire `start` on the non-zero case so this is a no-op when Syncthing
# is healthy. ExitOnFailure keeps systemd from flapping if the start fails;
# the next timer tick will retry.
ExecStart=/bin/bash -c "/usr/bin/systemctl is-active --quiet syncthing@%i.service || /usr/bin/systemctl start syncthing@%i.service"

# We intentionally do NOT set Restart= — a oneshot that succeeds is the
# whole point; the timer is what schedules the next probe.

[Install]
WantedBy=multi-user.target
